# modem_monitor.service
[Unit]
Description=Modem monitor

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/modem_monitor.sh

[Install]
WantedBy=multi-user.target
